# crontab
* * * * * ubuntu /bin/bash /home/ubuntu/apps/elections16graphics/repository/screenshotter/shot.sh > /var/log/elections16graphics/crontab.log 2>&1
