# crontab
* * * * * ubuntu /bin/bash /home/ubuntu/apps/elections16graphics/repository/screenshotter/shot.sh 2>&1
